# cd-tools/deepin-installer-for-dclc/deepin-installer/lib/systemd/system/systemd-timesyncd.service.d/disable-timesyncd-with-installer.conf
[Unit]
# Don't run timesyncd if installer is installed
# 1. 安装器阶段启动时间同步服务会导致无法设置时间
# 2. 安装器阶段NTP关闭后的一小段时间内（例如5~8秒），用户设置时间会失败(https://github.com/systemd/systemd/issues/11420)
# 3. 所以安装阶段不启动时间同步服务，安装器卸载该文件被移除后，时间同步服务则可正常启动
ConditionFileIsExecutable=!/usr/bin/deepin-installer
